# deploy/fishgoo-mcp/fishgoo-mcp.service
[Unit]
Description=Fishgoo Ads OS MCP Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/raddit
EnvironmentFile=/root/raddit/.env.fishgoo-mcp
ExecStart=/root/raddit/.venv/bin/python -m apps.fishgoo_mcp.server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
